# resources/systemd.service
[Unit]
Description=%NAME%

[Service]
ExecStart=yarn start %LISTEN%
WorkingDirectory=%CWD%
Restart=always

[Install]
WantedBy=default.target
